# xtrn-scripts-and-config/deb-file-creation/postinst
#!/bin/bash

# Install Java and tomcat
# apt update
# apt install default-jdk
useradd -r -m -U -d /opt/tomcat -s /bin/false tomcat
mkdir /opt/tomcat
tar xf /tmp/apache-tomcat-9*.tar.gz -C /opt/tomcat
ln -s /opt/tomcat/apache-tomcat-9.0.27 /opt/tomcat/latest

chown -RH tomcat: /opt/tomcat/latest
sh -c 'chmod +x /opt/tomcat/latest/bin/*.sh'
TOMCAT_SERVICE_FILE=/tmp/tomcat-service

cat $TOMCAT_SERVICE_FILE > /etc/systemd/system/tomcat.service

systemctl daemon-reload
systemctl enable tomcat
sudo usermod -a -G video tomcat

rm /tmp/apache-tomcat-9*.tar.gz
rm /tmp/server.xml
rm $TOMCAT_SERVICE_FILE